# configs/com.alexrabarts.focus-agent.ngrok.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexrabarts.focus-agent.ngrok</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ngrok</string>
        <string>start</string>
        <string>--config</string>
        <string>/Users/alex/.focus-agent/ngrok.yaml</string>
        <string>--all</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alex/.focus-agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/alex/.focus-agent/log/ngrok.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alex/.focus-agent/log/ngrok-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
